# templates/reaktor.service.erb
[Service]
Type=forking
PermissionsStartOnly=true
ExecStartPre=/usr/bin/mkdir -p /run/reaktor
ExecStartPre=/usr/bin/chown reaktor:reaktor /run/reaktor
ExecStart=/usr/bin/rake start
ExecStop=/usr/bin/rake stop
EnvironmentFile=<%= scope['::reaktor::homedir'] %>/etc/reaktor_environment
WorkingDirectory=<%= scope['::reaktor::_install_dir']%>
User=<%= scope['::reaktor::user'] %>
Group=<%= scope['::reaktor::group'] %>
PIDFile=/run/reaktor/reaktor.<%= scope['::reaktor::port'] %>.pid
